# systemd/ats-monitoring.service
[Unit]
Description=ATS 2.0 Monitoring Daemon
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=ats
Group=ats
WorkingDirectory=/opt/ats-system
Environment="PATH=/opt/ats-system/venv/bin:/usr/bin"
Environment="PYTHONPATH=/opt/ats-system"
EnvironmentFile=/opt/ats-system/.env

# Main process
ExecStart=/opt/ats-system/venv/bin/python /opt/ats-system/monitoring/websocket_daemon.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/ats/monitoring.log
StandardError=append:/var/log/ats/monitoring.error.log

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target